# dist/purenas.service
[Install]
WantedBy=multi-user.target

[Unit]
Description=PureNas autostart
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/etc/PureNAS/init